# packaging/systemd/ployz-runtime.service
[Unit]
Description=Ployz runtime reconciler
After=network-online.target docker.service ployzd.service
Wants=network-online.target docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ployz-runtime --data-root /var/lib/ployz/networks
Restart=always
RestartSec=2
User=root
Group=root

[Install]
WantedBy=multi-user.target
