# packaging/systemd/ployzd.service
[Unit]
Description=Ployz control plane daemon
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/ployzd
ExecStart=/usr/local/bin/ployzd --data-dir /var/lib/ployz run --mode host-service
Restart=always
RestartSec=2
User=root
Group=root

[Install]
WantedBy=multi-user.target
